# docker/supervisord.cfg
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
user=smartcinema

[program:api]
command=/usr/local/bin/smartcinema-api
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:web]
command=/usr/local/bin/smartcinema-web
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
